# zxhnh298n_hv14_fv113_hyperoptic_firmware/etc/inittab
::sysinit:/etc/rc
::respawn:/sbin/getty ttyS0 115200
